# .config/systemd/user/emacs-26.service.d/local.conf
[Service]
# added in order for having complete bash environment loaded
# with emacs including full PATH, etc.
ExecStart=
ExecStart=/bin/bash -l -c 'exec /usr/bin/emacs-26 --fg-daemon'
